# salt/setip.timer
[Unit]
Description=Update dynamic IPs in DNS zones (timer)

[Timer]
OnBootSec=60
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
